# deploy/systemd/octra-vitals-lab-history-mirror.timer
[Unit]
Description=Run Octra Vitals optional Lab history mirror

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
