# deploy/systemd/octra-vitals-notify-digest.timer
[Unit]
Description=Run Octra Vitals operator digest

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
