# deploy/systemd/octra-vitals-lab-history-mirror.timer
[Unit]
Description=Repair Octra Vitals optional Lab history mirror

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
