# deploy/systemd/octra-vitals-gateway.service
[Unit]
Description=Octra Vitals gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=octra-vitals
WorkingDirectory=/opt/octra-vitals/current
Environment=HOST=127.0.0.1
Environment=PORT=4173
Environment=VITALS_GATEWAY_ROLE=production
Environment=VITALS_DATA_DIR=/var/lib/octra-vitals
Environment=VITALS_STATE_SOURCE_MODE=program_required
Environment=VITALS_STATIC_ASSET_SOURCE=circle_required
Environment=VITALS_TRAFFIC_DIR=/var/lib/octra-vitals/traffic
EnvironmentFile=-/etc/octra-vitals/gateway.env
ExecStart=/usr/bin/node /opt/octra-vitals/current/dist/gateway/node-gateway/server.js
Restart=always
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/octra-vitals
RestrictSUIDSGID=yes
LockPersonality=yes
RemoveIPC=yes
CapabilityBoundingSet=
AmbientCapabilities=
ProtectProc=invisible
ProcSubset=pid
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
SystemCallArchitectures=native
UMask=0007

[Install]
WantedBy=multi-user.target
