# deploy/parking-page/octra-vitals-parking.service
[Unit]
Description=Octra Vitals parking page
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=octra-vitals
Group=octra-vitals
WorkingDirectory=/opt/octra-vitals-parking/current
ExecStart=/usr/bin/python3 -m http.server 8000 --bind 127.0.0.1
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/octra-vitals-parking
PrivateDevices=true
RestrictSUIDSGID=true
LockPersonality=true
RemoveIPC=true
ProtectKernelLogs=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
